# services/api-orchestrator/com.openclaw.api-orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.api-orchestrator</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/laurenz/.openclaw/workspace/services/api-orchestrator/server.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/laurenz/.openclaw/workspace/services/api-orchestrator</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>ORCHESTRATOR_PORT</key>
        <string>3005</string>
        <key>OPENCLAW_GATEWAY_URL</key>
        <string>http://127.0.0.1:8000</string>
        <key>OPENCLAW_GATEWAY_TOKEN</key>
        <string>***REDACTED***</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/laurenz/.openclaw/workspace/services/api-orchestrator/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/laurenz/.openclaw/workspace/services/api-orchestrator/stderr.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
