# services/gateway-proxy/com.openclaw.gateway-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.gateway-proxy</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/laurenz/.openclaw/workspace/services/gateway-proxy/server.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/laurenz/.openclaw/workspace/services/gateway-proxy</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>GATEWAY_PROXY_PORT</key>
        <string>8000</string>
        <key>REAL_GATEWAY_URL</key>
        <string>http://127.0.0.1:18789</string>
        <key>ORCHESTRATOR_URL</key>
        <string>http://127.0.0.1:3005</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/laurenz/.openclaw/workspace/services/gateway-proxy/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/laurenz/.openclaw/workspace/services/gateway-proxy/stderr.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
